# kiosk-url.service
[Unit]
Description=Restart kiosk browser on URL change

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart kiosk-browser.service